# pandabot.service
[Unit] 
Description=Panda Bot

[Service]
ExecStart=python3 /home/kstawik/pandabot/main.py

[Install] 
WantedBy=multi-user.target